# root/etc/services.d/watcher/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /app/watcher/watcher.py \
	--conf /config/config.cfg \
	--address 0.0.0.0 \
	--log /config/log \
	--port 9090 \
