# root/etc/services.d/watcher/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc python /app/watcher/watcher.py \
	--config_file=/config/config.cfg
